# email-report.timer
[Unit]
Description=Run Noise Detector Email Report every 2 hours
Requires=email-report.service

[Timer]
# Run every 2 hours
OnCalendar=*:0/2:00
# Also run 2 hours after boot
OnBootSec=2h
# Run 2 hours after the last run (if missed)
OnUnitActiveSec=2h
Persistent=true

[Install]
WantedBy=timers.target

