# noise-monitor.service.example
[Unit]
Description=Noise Monitor and Recorder
After=network.target sound.target

[Service]
WorkingDirectory=${SERVICE_WORKING_DIR}
ExecStart=/usr/bin/python3 ${SERVICE_WORKING_DIR}/noise_detector.py monitor
Restart=always
RestartSec=10
User=${SERVICE_USER}

# Capture stdout/stderr to systemd journal
StandardOutput=journal
StandardError=journal

# Environment variables (if needed)
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target

