# systemd/email-report.service.example
[Unit]
Description=Noise Detector Email Report
After=network.target

[Service]
Type=oneshot
WorkingDirectory=${SERVICE_WORKING_DIR}
ExecStart=/usr/bin/python3 ${SERVICE_WORKING_DIR}/scripts/email_report.py
User=${SERVICE_USER}

# Capture output to systemd journal
StandardOutput=journal
StandardError=journal

# Environment variables (can override config.json)
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target

